# systemd/pirecorder-webui.service
[Unit]
Description=pi-recorder: webui

[Service]
ExecStart=/usr/bin/node %h/pi-recorder/webui/server.js

[Install]
WantedBy=default.target
