# systemd/pirecorder-record.service
[Unit]
Description=pi-recorder: gpio-recorder.py

[Service]
ExecStart=%h/pi-recorder/gpio-recorder.py

[Install]
WantedBy=default.target
